# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:redis]
command=redis-server --port 6379 --daemonize no --appendonly no --maxmemory 100mb --maxmemory-policy allkeys-lru
autostart=true
autorestart=true
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis.log
priority=1

[program:api]
command=uvicorn src.webhook.app:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/api.log
stderr_logfile=/app/logs/api.log
priority=2

[program:worker]
command=python -m src.workers.extraction_worker
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/logs/worker.log
stderr_logfile=/app/logs/worker.log
priority=3
